# .config/systemd/user/battery-alert.service
[Unit]
Description=Battery Notification Service

[Service]
ExecStart=/usr/local/bin/battery-alertt

