# .config/systemd/user/charging-notification.service
[Unit]
Description=Battery Charging Notification Service

[Service]
ExecStart=/usr/local/bin/charging-notification

